# workloads/br-base/overlay/etc/inittab
::respawn:-/bin/bash /root/machine-boot-script.sh
